# home/Library/LaunchAgents/io.github.usami-k.brewupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTDPLIST1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.github.usami-k.brewupdate</string>
        <key>ProgramArguments</key>
        <array>
            <string>/opt/homebrew/bin/brew</string>
            <string>update</string>
        </array>
        <key>StandardOutPath</key>
        <string>/Users/usami/Library/Logs/io.github.usami-k.brewupdate.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/usami/Library/Logs/io.github.usami-k.brewupdate.log</string>
        <key>RunAtLoad</key>
        <true/>
        <key>StartInterval</key>
        <integer>3600</integer>
    </dict>
</plist>
