# home/Library/LaunchAgents/io.github.usami-k.podrepoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTDPLIST1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.github.usami-k.podrepoupdate</string>
        <key>ProgramArguments</key>
        <array>
            <string>/bin/sh</string>
            <string>-c</string>
            <string>echo ; date ; /Users/usami/.rbenv/shims/pod repo update</string>
        </array>
        <key>EnvironmentVariables</key>
        <dict>
            <key>LANG</key>
            <string>en_US.UTF-8</string>
        </dict>
        <key>StandardOutPath</key>
        <string>/Users/usami/Library/Logs/io.github.usami-k.podrepoupdate.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/usami/Library/Logs/io.github.usami-k.podrepoupdate.log</string>
        <key>RunAtLoad</key>
        <true/>
        <key>StartInterval</key>
        <integer>3600</integer>
    </dict>
</plist>
